# hermes-gateway-zenos.service
[Unit]
Description=Hermes Agent Gateway - Zenos Profile
After=network-online.target zenos-runtime.service srv-etla-workspaces.mount
Wants=network-online.target zenos-runtime.service srv-etla-workspaces.mount
StartLimitIntervalSec=0

[Service]
Type=simple
# The operator explicitly requires Hermes to own the complete VPS control plane.
# Keep this root-authoritative across restarts and Runtime redeployments.
User=root
Group=root
WorkingDirectory=/root
Environment=HOME=/root
Environment=USER=root
Environment=LOGNAME=root
Environment=PATH=/usr/local/lib/hermes-agent/venv/bin:/usr/local/lib/hermes-agent/node_modules/.bin:/usr/local/bin:/usr/bin:/bin
Environment=VIRTUAL_ENV=/usr/local/lib/hermes-agent/venv
Environment=PYTHONPATH=/usr/local/lib/hermes-agent/venv/lib/python3.11/site-packages:/usr/local/lib/hermes-agent
Environment=HERMES_HOME=/var/lib/hermes/.hermes/profiles/zenos
Environment=HERMES_CREDENTIAL_NAME=hermes-zenos.env
LoadCredentialEncrypted=hermes-zenos.env:/etc/credstore.encrypted/hermes-zenos.env.cred
ExecStart=/usr/local/lib/hermes-python/bin/python3.11 /usr/local/sbin/run-hermes-gateway-with-credential
Restart=always
RestartSec=5
RestartForceExitStatus=75
KillMode=mixed
KillSignal=SIGTERM
ExecReload=/bin/kill -USR1 $MAINPID
ExecStopPost=-/usr/local/lib/hermes-python/bin/python3.11 -m gateway.cgroup_cleanup
TimeoutStopSec=210
UMask=0077
# Root-authoritative host mode: do not hide /root, remount the host read-only,
# remove capabilities, isolate devices, or block kernel/service administration.
NoNewPrivileges=false
PrivateTmp=false
PrivateDevices=false
ProtectSystem=off
ProtectHome=false
ProtectHostname=false
ProtectClock=false
ProtectKernelTunables=false
ProtectKernelModules=false
ProtectKernelLogs=false
ProtectControlGroups=false
ProtectProc=default
ProcSubset=all
RestrictNamespaces=false
RestrictSUIDSGID=false
LockPersonality=false
RestrictRealtime=false
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
