# hermes-gateway-zenos.service
[Unit]
Description=Hermes Agent Gateway - Zenos Profile
After=network-online.target zenos-runtime.service zenos-memory.service
Wants=network-online.target zenos-runtime.service zenos-memory.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/usr/local/lib/hermes-agent
Environment=HOME=/var/lib/hermes
Environment=USER=root
Environment=LOGNAME=root
Environment=PATH=/usr/local/lib/hermes-agent/venv/bin:/usr/local/lib/hermes-agent/node_modules/.bin:/usr/local/bin:/usr/bin:/bin
Environment=VIRTUAL_ENV=/usr/local/lib/hermes-agent/venv
Environment=HERMES_HOME=/var/lib/hermes/.hermes/profiles/zenos
Environment=HERMES_CREDENTIAL_NAME=hermes-zenos.env
LoadCredentialEncrypted=hermes-zenos.env:/etc/credstore.encrypted/hermes-zenos.env.cred
ExecStart=/usr/local/lib/hermes-agent/venv/bin/python /opt/zenos-runtime/current/scripts/run-hermes-gateway-with-credential.py
Restart=always
RestartSec=5
RestartForceExitStatus=75
KillMode=mixed
KillSignal=SIGTERM
ExecReload=/bin/kill -USR1 $MAINPID
ExecStopPost=-/usr/local/lib/hermes-agent/venv/bin/python -m gateway.cgroup_cleanup
TimeoutStopSec=210
UMask=0077
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
