# zenos-runtime-backup.service
[Unit]
Description=Encrypted verified backup for Zenos Runtime
After=zenos-runtime.service
Requires=zenos-runtime.service

[Service]
Type=oneshot
User=zenos-runtime
Group=zenos-runtime
WorkingDirectory=/opt/zenos-runtime/current
LoadCredentialEncrypted=zenos-runtime.env:/etc/credstore.encrypted/zenos-runtime.env.cred
Environment=HOME=/var/lib/zenos-runtime
Environment=NODE_ENV=production
Environment=ZENOS_RUNTIME_DB_PATH=/var/lib/zenos-runtime/runtime.db
Environment=ZENOS_RUNTIME_BACKUP_DIR=/var/backups/zenos-runtime
Environment=ZENOS_RUNTIME_BACKUP_KEEP=14
ExecStart=/usr/bin/node /opt/zenos-runtime/current/scripts/backup-zenos-runtime.mjs
UMask=0077
NoNewPrivileges=true
CapabilityBoundingSet=
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadOnlyPaths=/opt/zenos-runtime/current /etc/zenos-runtime /var/lib/zenos-runtime
ReadWritePaths=/var/backups/zenos-runtime
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
TimeoutStartSec=120
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zenos-runtime-backup
